# server/deploy/xiaozhi-family-hub.service
[Unit]
Description=Xiaozhi Family Hub Backend
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/xiaozhi-family-hub/current
Environment=HOST=0.0.0.0
Environment=PORT=3100
Environment=PUBLIC_BASE_URL=https://wave.xionglei.online
Environment=DATA_DIR=/opt/xiaozhi-family-hub/shared/data
Environment=STORE_DRIVER=sqlite
Environment=SQLITE_FILE=/opt/xiaozhi-family-hub/shared/data/family-hub.sqlite
Environment=SQLITE_RELATIONAL_READS=1
Environment=SQLITE_COMPACT_KV=1
Environment=RESOURCE_DIR=/opt/xiaozhi-family-hub/shared/resources
Environment=SERVER_MUSIC_DIR=/opt/xiaozhi-family-hub/shared/resources/music/server
Environment=ADMIN_PUBLIC_HOSTS=wave.xionglei.online
Environment=ADMIN_PROTECT_PUBLIC_MUTATIONS=1
Environment=ALLOW_QUERY_TOKENS=0
Environment=CORS_ALLOWED_ORIGINS=https://wave.xionglei.online,http://192.168.31.246:3100
Environment=WEATHER_PROVIDER=open-meteo
Environment=ONLINE_MEDIA_PROVIDERS=rss,archive,direct
Environment=ONLINE_MEDIA_FEEDS=
Environment=ONLINE_MEDIA_ALLOW_DIRECT_URLS=1
Environment=ONLINE_MEDIA_ALLOW_PRIVATE_HOSTS=0
Environment=ONLINE_MEDIA_TIMEOUT_MS=4500
Environment=ONLINE_MEDIA_MAX_RESULTS=6
Environment=NAS_MUSIC_DIR=
Environment=OPENCLAW_COMMAND=/opt/xiaozhi-family-hub/current/bin/openclaw-command-adapter.sh
Environment=OPENCLAW_TARGETS=default,music,diagnostics
Environment=HOME_ASSISTANT_URL=
Environment=HOME_ASSISTANT_TOKEN=
Environment=HOME_ASSISTANT_SCENES=scene.family,scene.home,scene.evening,scene.good_night,scene.morning,scene.away
Environment=HOME_ASSISTANT_TIMEOUT_MS=5000
EnvironmentFile=-/opt/xiaozhi-family-hub/shared/secrets/admin.env
EnvironmentFile=-/opt/xiaozhi-family-hub/shared/secrets/integrations.env
EnvironmentFile=-/opt/xiaozhi-family-hub/shared/secrets/home-assistant.env
ExecStart=/opt/node-v22/bin/node index.js
Restart=always
RestartSec=3
User=xiaozhi
Group=xiaozhi

[Install]
WantedBy=multi-user.target
